# misc/profile.d/deepin-deb-fix.sh
#!/bin/bash

export XDG_DATA_DIRS="/usr/share/deepin-desktop-fix:$XDG_DATA_DIRS"